# DEBIAN/foxtrot-server/postrm
#! /bin/bash
set -e

PKG_NAME=foxtrot-server
update-rc.d -f $PKG_NAME remove

exit 0
